# deploy/sentinel-dns.service
[Unit]
Description=Sentinel DNS — RAM-first ad-blocking DNS resolver
After=network-online.target docker.service
Wants=network-online.target
Documentation=https://github.com/Alexanderrrrrrw/SentinelDNS

[Service]
Type=simple
WorkingDirectory=/opt/sentinel-dns/deploy
EnvironmentFile=-/opt/sentinel-dns/deploy/.env
ExecStart=docker compose up --no-build
ExecStop=docker compose down

# SD-card wear leveling: idle I/O class so checkpoint flushes
# never compete with DNS packet handling.
IOSchedulingClass=idle
IOSchedulingPriority=7

Nice=5
CPUSchedulingPolicy=batch

# Prefer killing Sentinel over the kernel or sshd on OOM.
OOMScoreAdjust=200

Restart=always
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
